# pi-goes-ns/systemd/goesrecv-watchdog.timer
[Unit]
Description=Run goesrecv SDR-stall watchdog every 2 min
[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=20s
[Install]
WantedBy=timers.target
