# lxc-goes-proc/systemd/goes-prune.timer
[Unit]
Description=Hourly prune of local GOES cache
[Timer]
OnBootSec=30min
OnUnitActiveSec=1h
[Install]
WantedBy=timers.target
